# contrib/packages/archlinux/nodejs-galilel-explorer-sync-rich.timer
[Unit]
Description = Galilel Explorer (Sync) - Rich list update timer

[Timer]
OnCalendar = *:0/1
Persistent = true
 
[Install]
WantedBy = timers.target
