# contrib/packages/archlinux/nodejs-galilel-explorer-sync-peer.service
[Unit]
Description = Galilel Explorer (Sync) - Network peers update

[Service]
Type = simple
WorkingDirectory = /usr/lib/node_modules/galilel-explorer
ExecStart = /usr/bin/node cron/peer.js
SyslogIdentifier = galilel-explorer
StandardOutput = syslog
StandardError = syslog
