# contrib/packages/archlinux/nodejs-galilel-explorer-sync-proposal.service
[Unit]
Description = Galilel Explorer (Sync) - Proposal update

[Service]
Type = simple
WorkingDirectory = /usr/lib/node_modules/galilel-explorer
ExecStart = /usr/bin/node cron/proposal.js
SyslogIdentifier = galilel-explorer
StandardOutput = syslog
StandardError = syslog
