# contrib/packages/archlinux/nodejs-galilel-explorer-sync-block.timer
[Unit]
Description = Galilel Explorer (Sync) - Blockchain update timer

[Timer]
OnCalendar = *:0/1
Persistent = true
 
[Install]
WantedBy = timers.target
