# contrib/packages/archlinux/nodejs-galilel-explorer-sync-block.service
[Unit]
Description = Galilel Explorer (Sync) - Blockchain update

[Service]
Type = simple
WorkingDirectory = /usr/lib/node_modules/galilel-explorer
ExecStart = /usr/bin/node cron/block.js
SyslogIdentifier = galilel-explorer
StandardOutput = syslog
StandardError = syslog
