# contrib/packages/archlinux/nodejs-galilel-explorer-sync-masternode.service
[Unit]
Description = Galilel Explorer (Sync) - Masternode update

[Service]
Type = simple
WorkingDirectory = /usr/lib/node_modules/galilel-explorer
ExecStart = /usr/bin/node cron/masternode.js
SyslogIdentifier = galilel-explorer
StandardOutput = syslog
StandardError = syslog
